# debian/postinst
#!/bin/sh
# postinst script for fuckooredoo

set -e

case "$1" in
    configure)
        # Ensure Tor service is enabled
        if command -v systemctl >/dev/null 2>&1; then
            if systemctl list-unit-files | grep -q tor.service; then
                echo "Enabling Tor service..."
                systemctl enable tor.service || true
                
                # Start Tor if not running
                if ! systemctl is-active --quiet tor.service; then
                    echo "Starting Tor service..."
                    systemctl start tor.service || true
                fi
            fi
        fi
        
        # Update desktop database for icon and menu integration
        if command -v update-desktop-database >/dev/null 2>&1; then
            update-desktop-database -q /usr/share/applications || true
        fi
        
        # Update icon cache
        if command -v gtk-update-icon-cache >/dev/null 2>&1; then
            gtk-update-icon-cache -q -t -f /usr/share/icons/hicolor || true
        fi
        
        echo ""
        echo "â•”â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•—"
        echo "â•‘     ðŸ–• F**k Ooredoo - Installation Complete! ðŸ–•          â•‘"
        echo "â•šâ•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•â•"
        echo ""
        echo "âœ… Package installed successfully!"
        echo ""
        echo "ðŸš€ To launch the application, run:"
        echo "   $ fuckooredoo"
        echo ""
        echo "ðŸ“š For help and documentation:"
        echo "   $ fuckooredoo --help"
        echo "   $ man fuckooredoo"
        echo ""
        echo "ðŸ”§ Tor service status:"
        if systemctl is-active --quiet tor.service 2>/dev/null; then
            echo "   âœ… Tor is running (ready to bypass ISP throttling!)"
        else
            echo "   âš ï¸  Tor is not running. Start it with:"
            echo "      $ sudo systemctl start tor"
        fi
        echo ""
        echo "ðŸ“– Documentation: /usr/share/doc/fuckooredoo/"
        echo "ðŸŒ Homepage: https://github.com/AzizBahloul/F_ooredoo"
        echo ""
        echo "ðŸŽ¯ Ready to bypass ISP throttling! Fuck Ooredoo! ðŸŽ¯"
        echo ""
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

#DEBHELPER#

exit 0
